# scripts/beaglebone_black/ttyO4.conf
# /etc/init/ttyO4.conf
start on startup

exec /home/ubuntu/kks/ttyO4.sh


